# services/shutdownbutton.service
[Unit]
Description=Shutdown Daemon

[Service]
ExecStart=/usr/local/sbin/eiclab_shutdown_button.py
Type=simple

[Install]
WantedBy=multi-user.target
